# service/nm-applet.service
[Unit]
Description=nm-applet

[Service]
ExecStart=nm-applet

[Install]
WantedBy=default.target
